# scripts/debian/postinst
if [ ! -f /usr/local/nginx/conf/passwd ]; then \
 touch /usr/local/nginx/conf/passwd;\
fi
chown -R http:http /usr/local/nginx
chmod 550  /usr/local/nginx/
chmod 550  /usr/local/nginx/conf/
chmod 440  /usr/local/nginx/conf/*
if /bin/systemctl >/dev/null 2>&1; then \
 /usr/sbin/update-rc.d nginx disable || true; \
 /usr/sbin/chkconfig nginx off || true; \
 /bin/systemctl enable nginx.service; \
elif [ -x /usr/sbin/update-rc.d ]; then \
 /usr/sbin/update-rc.d nginx defaults; \
 /usr/sbin/update-rc.d nginx enable; \
elif [ -x /usr/sbin/chkconfig ]; then \
 /usr/sbin/chkconfig nginx on; \
fi

echo ""
echo nginx installed.
if /bin/systemctl >/dev/null 2>&1; then \
 /bin/systemctl stop nginx.service \
 /bin/systemctl start nginx.service \
 echo use ;\
 echo "  "systemctl start nginx.service ;\
 echo to start http ;\
 echo "  "systemctl stop nginx.service ;\
 echo to stop http ;\
elif [ -x /usr/sbin/service ]; then \
 /usr/sbin/service nginx stop  || true;\
 /usr/sbin/service nginx start  || true;\
 echo "  "service nginx start ;\
 echo to start http ;\
 echo "  "service nginx stop ;\
 echo to stop http ;\
fi
echo "  "/usr/local/nginx/conf/addnginxuser.sh
echo to add users
echo ""
echo Default config uses Google\'s DNS.
echo It\'s recommended to use provider supplied DNS or install local recursor, e.g. pdns-recursor.
echo Configure preferred DNS in /usr/local/nginx/conf/nginx.cfg.
echo run \'/usr/local/nginx/conf/addnginxuser.sh admin password\' to configure \'admin\' user
